# ops/cicdbuzz.conf
[program:cicdbuzz]
directory=/opt/cicd-buzz
command=python3 app.py
autostart=true
autorestart=true